# samp.service
[Unit]
Description=Greatbone Sample Service

[Service]
EnvironmentFile=/etc/environment
Type=forking
WorkingDirectory=/root
ExecStart=/usr/bin/dotnet /root/publish/Greatbone.dll deploy
Restart=always

[Install]
WantedBy=multi-user.target
